# docker/supervisor/configs/Client.Authentication.Application.Commands.Login.ini

[program:Shop.Client.Authentication.Application.Commands.Login]
process_name    = %(process_num)02d
command         = php /var/project/artisan rabbitmq:command-consume --queue Client.Authentication.Application.Commands.Login
autostart       = true
autorestart     = true
user            = root
numprocs        = 1
redirect_stderr = true
stdout_logfile  = /var/project/storage/logs/Client.Authentication.Application.Commands.Login.log
