# docker/supervisor/configs/commands-and-domain-events/Admin.Currency.Application.Commands.Create.ini

[program:Shop.Admin.Currency.Application.Commands.Create]
process_name    = %(process_num)02d
command         = php /var/project/artisan rabbitmq:command-consume --queue Admin.Currency.Application.Commands.Create
autostart       = true
autorestart     = true
user            = root
numprocs        = 1
redirect_stderr = true
stdout_logfile  = /var/project/storage/logs/Admin.Currency.Application.Commands.Create.log
