# docker/supervisor/configs/commands-and-domain-events/Client.Authentication.Application.Commands.RestorePasswordCode.ini

[program:Shop.Client.Authentication.Application.Commands.RestorePasswordCode]
process_name    = %(process_num)02d
command         = php /var/project/artisan rabbitmq:command-consume --queue Client.Authentication.Application.Commands.RestorePasswordCode
autostart       = true
autorestart     = true
user            = root
numprocs        = 1
redirect_stderr = true
stdout_logfile  = /var/project/storage/logs/Client.Authentication.Application.Commands.RestorePasswordCode.log
