# docker/supervisor/configs/commands-and-domain-events/Client.Address.Application.Commands.Create.ini

[program:Shop.Client.Address.Application.Commands.Create]
process_name    = %(process_num)02d
command         = php /var/project/artisan rabbitmq:command-consume --queue Client.Address.Application.Commands.Create
autostart       = true
autorestart     = true
user            = root
numprocs        = 1
redirect_stderr = true
stdout_logfile  = /var/project/storage/logs/Client.Address.Application.Commands.Create.log
