# docker/supervisor/configs/commands-and-domain-events/Client.Favorite.Application.Subscribers.MemberWasRegisteredDomainEventSubscriber.ini

[program:Shop.Client.Favorite.Application.Subscribers.MemberWasRegisteredDomainEventSubscriber]
process_name    = %(process_num)02d
command         = php /var/project/artisan rabbitmq:domain-event-consume --queue Client.Favorite.Application.Subscribers.MemberWasRegisteredDomainEventSubscriber
autostart       = true
autorestart     = true
user            = root
numprocs        = 1
redirect_stderr = true
stdout_logfile  = /var/project/storage/logs/Client.Favorite.Application.Subscribers.MemberWasRegisteredDomainEventSubscriber.log
